# deploy/portfolio.service
[Unit]
Description=Portfolio Website
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=portfolio
Group=portfolio
WorkingDirectory=/var/www/portfolio
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=3000
Environment=BODY_SIZE_LIMIT=8M
Environment=SHUTDOWN_TIMEOUT=30
Environment=IDLE_TIMEOUT=60
EnvironmentFile=/etc/portfolio/portfolio.env
ExecStart=/usr/bin/node build/index.js
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=45
StateDirectory=portfolio
LogsDirectory=portfolio
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/portfolio /var/www/portfolio/static/uploads /var/www/portfolio/static/assets/work
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
AmbientCapabilities=
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target
